# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:novnc]
command=/root/noVNC/utils/launch.sh --vnc localhost:5900 --listen 8080
autorestart=true
